# server/src/unix-control/deb/postinst
#!/bin/sh
#

echo "Registering ${unix.service} service..."
update-rc.d ${unix.service} defaults
